# Shell Script/Jamf to Workspace One Migration PKG/payload/Library/LaunchDaemons/com.vmware.migrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.vmware.migrator</string>
<key>ProgramArguments</key>
<array>
	<string>/bin/bash</string>
	<string>/Library/Application Support/VMware/migrator.sh</string>
	<string>--origin</string>
	<string>custom</string>
	<string>--registration-type</string>
	<string>prompt</string>
	<string>--dest-baseurl</string>
	<string>your-ws1-server-url</string>
	<string>--dest-auth</string>
	<string>Basic your-ws1-api-id</string>
	<string>--dest-token</string>
	<string>your-ws1-auth-token</string>
	<string>--dest-groupid</string>
	<string>your-ws1-group/location-id-integer</string>
	<string>--dest-apiurl</string>
	<string>your-ws1-api-url</string>
	<string>--user-prompt</string>
	<string>email</string>
	<string>--removal-script</string>
	<string>/Library/Application Support/VMware/MigratorResources/removemdm.sh</string>
</array>
<key>RunAtLoad</key>
<true/>
<key>StandardErrorPath</key>
<string>/var/log/vmw_migrator.log</string>
<key>StandardOutPath</key>
<string>/var/log/vmw_migrator.log</string>
</dict>
</plist>
